# web_ui/frontend/dev/image/crontab
*/28 * * * * /opt/bin/login.sh > /proc/1/fd/1
